# shadowsocks-rust_server.service@openwrt
#!/bin/sh /etc/rc.common
START=95

conf_folder="/etc/shadowsocks-rust_server"
pid_folder="/var/run/shadowsocks-rust_server"
log_file="/var/log/shadowsocks-rust_server.log"

do_start() {
	echo "开启ssserver进程..."  >> $log_file
	[ -d $conf_folder/ ] || return 2
	json_files=$(ls $conf_folder/ | grep ".json$" )
	[ -d $pid_folder/ ] && do_stop
	mkdir -p $pid_folder/
	#ulimit -n 51200
	process_num=1
	for json_file in $json_files; do
		pre_name=$(echo $json_file | sed "s/\.json//g")
		nohup ssserver -c $conf_folder/$json_file  > /dev/null 2>&1 &
		echo $! > $pid_folder/$pre_name.pid
	done
}

do_stop() {
	echo "关闭ssserver进程..."  >> $log_file
	[ -d $pid_folder/ ] || return 2
	pid_files=$(ls $pid_folder/ | grep ".pid" )
	[ -z "$pid_files" ] && return 2
	for pid_file in $pid_files; do
		pid=$(cat $pid_folder/$pid_file)
		kill $pid
	done
	rm -rf $pid_folder/
}

start() {
	rm -f $log_file
	do_start
}

stop() {
	rm -f $log_file
	do_stop
}

restart() {
	rm -f $log_file
	do_stop
	do_start
}
